# roles/openbox/templates/etc/systemd/system/getty@tty1.service.d/autologin.conf.j2
# {{ ansible_managed }}
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin {{ default_user.name }} --noclear %I 38400 linux

